# offlineimap/.config/systemd/user/offlineimap-oneshot.service.d/override.conf
[Service]
Environment=PASSWORD_STORE_DIR=/home/sfackler/Passwörter
Environment=GNUPGHOME=/home/sfackler/.config/gnupg
Environment=NOTMUCH_CONFIG=/home/sfackler/.config/notmuch/config
